# bot/bin/olivka_daemon_start
#!/bin/bash
### BEGIN INIT INFO
# Provides:          olivka_daemon_start
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Run olivka_daemon
### END INIT INFO

LOG=/home/olivka/bot/bin/logs/startup.log

echo "Starting OLIVKA daemon..." | tee -a $LOG
su olivka -c "nice -n 19 nohup /home/olivka/bot/bin/olivka_daemon 2>&1 | tee -a $LOG &" olivka

